# usr/systemd/system/slimlock.service
[Unit]
Description=Lock X session using slimlock

[Service]
User=gedi
ExecStart=/usr/bin/slimlock

[Install]
WantedBy=sleep.target
