# usr/systemd/system/xcompmgr.service
[Unit]
Description=Xcompmgr
Before=after-wm.target
After=wm.target

[Service]
ExecStart=/usr/bin/xcompmgr -f

[Install]
WantedBy=after-wm.target
